# machines/EXAMPLE-HOST/systemd/mnt-data.automount
[Unit]
Description=Automount for /mnt/data

[Automount]
Where=/mnt/data

[Install]
WantedBy=multi-user.target
